# Python-Bolsa/RaspberryPi/LaserSwitch/crontab
*/1 * * * * pi /home/pi/sw/laserSwitch.py
